# incus-cn-blocker-update.service
[Unit]
Description=Refresh China CIDR cache for Incus destination blocking
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/incus-cn-blocker update --apply-if-enabled
