# incus-cn-blocker-apply.service
[Unit]
Description=Apply Incus China destination block rules from cache
After=network-online.target nftables.service
Wants=network-online.target
ConditionPathExists=/var/lib/incus-cn-blocker/state.env

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/incus-cn-blocker apply-cached --if-enabled

[Install]
WantedBy=multi-user.target
